# pixie-attacker.service
[Unit]
Description=Pixie-Dust WPS Attacker Service
After=network.target multi-user.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/usr/local/bin/pixie-attacker-daemon.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
TimeoutStartSec=0
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
